# discord-music-bot.service
# systemd service file for the Discord Music Bot
#
# Installation:
#   1. Copy this file:      sudo cp discord-music-bot.service /etc/systemd/system/
#   2. Edit the file:       sudo systemctl edit --full discord-music-bot
#      - Set WorkingDirectory, User, EnvironmentFile, ExecStart, and ReadWritePaths
#   3. Reload systemd:      sudo systemctl daemon-reload
#   4. Enable on boot:      sudo systemctl enable discord-music-bot
#   5. Start the service:   sudo systemctl start discord-music-bot
#
# Logs:  journalctl -u discord-music-bot -f

[Unit]
Description=Discord Music Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# --- EDIT THESE TO MATCH YOUR SETUP ---
User=your-username
WorkingDirectory=/path/to/DiscordMusicPlayer
EnvironmentFile=/path/to/DiscordMusicPlayer/.env
ExecStart=/path/to/DiscordMusicPlayer/.venv/bin/discord-music-player
ReadWritePaths=/path/to/DiscordMusicPlayer/data /path/to/DiscordMusicPlayer/logs
# ---------------------------------------

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# Graceful shutdown (SIGTERM → 30s timeout → SIGKILL)
TimeoutStopSec=30
KillSignal=SIGTERM

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
